# pihole/macvlan.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/macvlan.sh

[Install]
WantedBy=default.target
